# packaging/debian/usc_package_xenial/debian/usc.prerm
#!/bin/bash -e

# Source debconf library.
. /usr/share/debconf/confmodule
set -e

CONF_PATH="/etc/usc"

systemctl disable usc.service

#DEBHELPER#
